# db/supervisor.conf
[program:mysqld]
command=/usr/sbin/mysqld
autorestart=false
